# deploy/daily-stoic.service
[Unit]
Description=Daily Stoic Website
Documentation=https://github.com/yourusername/daily-stoic-website
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/daily-stoic
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=daily-stoic
Environment=NODE_ENV=production
Environment=PORT=3000

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/daily-stoic

[Install]
WantedBy=multi-user.target 